# com.skatehive.cookie-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.skatehive.cookie-monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/vladnikolaev/skatehive-monorepo/skatehive-instagram-downloader/cookie-health-check.sh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>21600</integer> <!-- 6 hours in seconds -->
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/cookie-monitor.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/cookie-monitor.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>SKATEHIVE_MONOREPO</key>
        <string>/Users/vladnikolaev/skatehive-monorepo</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
